# .config/systemd/user/batteryListener.service
[Unit]
Description=Battery level listener for notifications
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/battery_listener.sh
Restart=always
StandardOutput=journal
StandardError=journal
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=%t/%u
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus

[Install]
WantedBy=graphical-session.target
